# sipsvs.conf
[program:sipsvs]
command=/home/ubuntu/sips-vs/main.py
user=ubuntu
autostart=true
autorestart=true
directory=/home/ubuntu/sips-vs/
